# files/postgresql-template.service
.include /lib/systemd/system/postgresql.service
[Unit]
Description=PostgreSQL Cluster %i

[Service]
Environment=PGDATA=/var/lib/pgsql/data/%I
